# linux/autostart/MyDotfiles.desktop
[Desktop Entry]
Name=Script
Type=Application
Exec=/usr/local/share/onlogin.sh
Terminal=true
StartupNotify=true
Icon=gnome-do
Hidden=false
